# com.mateo.digest-autodeploy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mateo.digest-autodeploy</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/mateomurray/Desktop/Digests/Market Analysis/auto-deploy-watch.sh</string>
  </array>
  <key>WatchPaths</key>
  <array>
    <string>/Users/mateomurray/Desktop/Digests/Market Analysis/semanales</string>
  </array>
  <key>StandardOutPath</key>
  <string>/Users/mateomurray/Desktop/Digests/Market Analysis/.launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/mateomurray/Desktop/Digests/Market Analysis/.launchd.err.log</string>
</dict>
</plist>
